# openclaw-workspaces/ai-order/launchd/com.ai-order.phase3-maintenance.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ai-order.phase3-maintenance</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/jinqianfei/openclaw-workspaces/ai-order/scripts/phase3_maintenance.sh</string>
    </array>

    <!-- 每周日 03:00 触发 -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>0</integer>
        <key>Hour</key>
        <integer>3</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/tmp/ai-order.phase3-maintenance.out.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/ai-order.phase3-maintenance.err.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/jinqianfei/openclaw-workspaces/ai-order</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>ExitTimeOut</key>
    <integer>600</integer>
</dict>
</plist>
